# media-server.service
[Unit]
Description=Local Media Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=media-server
Group=media-server
WorkingDirectory=/opt/media-server
ExecStart=/opt/media-server/.venv/bin/python /opt/media-server/server.py
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
